# board/r36s/rootfs-overlay/etc/inittab
# /etc/inittab — R36S CyberDeck OS minimal rootfs (Fase 2)
# BusyBox init. Objetivo: provar boot do nosso rootfs e abrir shell de debug.

# Montagens iniciais (proc/sys/dev/...) e banner:
::sysinit:/etc/init.d/rcS

# Shell no console serial de debug (ttyFIQ0) e na tela (tty1 = framebuffer).
# Confirmados funcionando no R36S (2026-06-15). Sem o respawn genérico sem tty
# (ele gerava "can't access tty; job control turned off").
ttyFIQ0::respawn:-/bin/sh
tty1::respawn:-/bin/sh

# Encerramento limpo:
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::restart:/sbin/init
