# board/r36s/rootfs-overlay/etc/inittab
# /etc/inittab — R36S CyberDeck OS minimal rootfs (Fase 2)
# BusyBox init. Objetivo: provar boot do nosso rootfs e abrir shell de debug.

# Montagens iniciais (proc/sys/dev/...) e banner:
::sysinit:/etc/init.d/rcS

# Serial de debug (ttyFIQ0): shell. Tela (tty1): a UI CyberDeck (renderizador de
# framebuffer); se ela sair/faltar, o launcher cai num shell na própria tela.
ttyFIQ0::respawn:-/bin/sh
tty1::respawn:/usr/local/bin/cyberdeck-launch.sh

# Encerramento limpo:
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::restart:/sbin/init
