# runtime/services/cyberdeck-agent.service
[Unit]
Description=CyberDeck data agent (Node.js — JSON de hardware/SO em 127.0.0.1:8080)
DefaultDependencies=no
After=systemd-udevd.service
Before=cyberdeck-x.service
Wants=cyberdeck-x.service

[Service]
Type=simple
# Debian pode instalar /usr/bin/nodejs em vez de /usr/bin/node — acha o que existir.
ExecStart=/bin/sh -c 'exec "$(command -v node || command -v nodejs)" /usr/local/lib/cyberdeck-agent/agent.js 8080'
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
