# runtime/services/cyberdeck-session.service
# systemd unit — sessão do CyberDeck: SELETOR de interface no boot + lança a UI escolhida.
# Substitui o cyberdeck-x.service como entrada de UI (a web continua sendo iniciada por
# start-cyberdeck-x.sh, agora chamado pelo cyberdeck-session.sh). Instalar em
# /etc/systemd/system/ (build-x11-rootfs.sh habilita esta unit no lugar de cyberdeck-x).

[Unit]
Description=R36S CyberDeck session (seletor de interface + UI)
After=systemd-user-sessions.service cyberdeck-agent.service
Wants=cyberdeck-agent.service
Conflicts=getty@tty1.service
After=getty@tty1.service

[Service]
Type=simple
# espera o /dev/fb0 (o seletor desenha nele)
ExecStartPre=/bin/sh -c 'for i in $(seq 1 20); do [ -e /dev/fb0 ] && exit 0; sleep 1; done; exit 0'
ExecStart=/usr/local/bin/cyberdeck-session.sh
Restart=always
RestartSec=2
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
